# monitoring/crontab
02 09 * * 1-5 /root/monitoring/start_monitoring.sh
02 22 * * 1-5 /root/monitoring/archive_monitoring.sh

